# debian/pi-usb-safegate/DEBIAN/prerm
#!/bin/bash
set -e

case "$1" in
    remove|deconfigure)
        if command -v systemctl >/dev/null 2>&1; then
            systemctl stop pi-usb-safegate.service || true
            systemctl disable pi-usb-safegate.service || true
        fi
        ;;
esac

exit 0
